# packaging/unix/com.croatiasecurity.behavedr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.croatiasecurity.behavedr</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/behavedr/Behavedr</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/behavedr</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>1</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/opt/behavedr/logs/behavedr-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/opt/behavedr/logs/behavedr-stderr.log</string>

    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>8192</integer>
    </dict>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DOTNET_ENVIRONMENT</key>
        <string>Production</string>
    </dict>

    <!-- v0.1.6 audit fixes (RT-10): Hardening keys -->

    <!-- Create isolated security session for the daemon -->
    <key>SessionCreate</key>
    <true/>

    <!-- Restrictive file creation mask (owner-only) -->
    <key>Umask</key>
    <integer>77</integer>

    <!-- Fast exit timeout for rapid re-launch on unload -->
    <key>ExitTimeOut</key>
    <integer>5</integer>

    <!-- Higher scheduling priority for security-critical monitoring -->
    <key>Nice</key>
    <integer>-10</integer>

    <!-- Do not deprioritize I/O — detection must be responsive -->
    <key>LowPriorityBackgroundIO</key>
    <false/>

    <!-- Prevent macOS from killing daemon under memory pressure -->
    <key>EnablePressuredExit</key>
    <false/>

    <!-- Register Mach service for IPC (future: admin CLI communication) -->
    <key>MachServices</key>
    <dict>
        <key>com.croatiasecurity.behavedr</key>
        <true/>
    </dict>
</dict>
</plist>
